# docker/satis/conf/supervisor/supervisor.conf
[supervisord]
nodaemon = true

[program:nginx]
user = root
command=/usr/sbin/nginx
autostart=true
autorestart=true
stdout_events_enabled=true